# deploy/systemd/systemecom-frontend.service
[Unit]
Description=SystemECOM Frontend (Next.js)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/systemecom/wlms-frontend
EnvironmentFile=/etc/systemecom/frontend.env
ExecStart=/usr/bin/npm run start -- --port 2020
Restart=on-failure
RestartSec=3

# Recommended (uncomment after creating a dedicated user + ensuring permissions):
# User=systemecom
# Group=systemecom
# UMask=0027

[Install]
WantedBy=multi-user.target


